# ops/systemd/nginx.service.d/override.conf
[Unit]
After=network-online.target

[Service]
LimitNOFILE=65535
Restart=on-failure
RestartSec=5
