# ops/systemd/hypercorn-codeintel.service
[Unit]
Description=CodeIntel MCP via Hypercorn
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=codeintel
Group=codeintel
WorkingDirectory=/opt/codeintel_rev
Environment=MCP_EDGE_MODE=edge-terminate
Environment=PROXY_TRUSTED_HOPS=1
Environment=PYTHONUNBUFFERED=1
Environment=PATH=/opt/codeintel_rev/.venv/bin:/usr/bin
EnvironmentFile=-/etc/codeintel_rev.env
ExecStart=/usr/bin/env bash -lc 'python -m hypercorn --config ops/hypercorn.toml codeintel_rev.app.main:asgi'
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=2
LimitNOFILE=65535
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
TimeoutStopSec=45
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target
