# config/systemd/search-api.service
[Unit]\nDescription=kgfoundry API\n[Service]\nExecStart=/usr/bin/uvicorn kgfoundry.search_api.app:app --app-dir src --port 8080\n