# config/systemd/search-api.service
[Unit]\nDescription=KGForge API\n[Service]\nExecStart=/usr/bin/uvicorn kgforge.search_api.app:app --app-dir src --port 8080\n