# config/systemd/vllm-embed.service
[Unit]\nDescription=vLLM Qwen3-Embedding-4B\n[Service]\nExecStart=/usr/bin/vllm serve ...\n
