# config/systemd/vllm-vlm.service
[Unit]\nDescription=vLLM Granite-Docling\n[Service]\nExecStart=/usr/bin/vllm serve ...\n
